# dcos-neo4j-proxy/nginx.crontab
* * * * * root /init-dcos-neo4j-nginx.sh

